# bastion/templates/systemd/edge-watchdog.service
[Unit]
Description=edge-watchdog — standing egress/relay self-heal
Documentation=man:nft(8)
After=network.target nftables.service wg-quick@{{ interfaces.wg_vps_iface }}.service
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/local/sbin/edge-watchdog
Restart=always
RestartSec=10
# F1 — a self-heal tool, so confinement is DELIBERATELY LIGHT: no ProtectSystem=strict and no
# SystemCallFilter (it must shell out to systemctl/ping/conntrack/net-rollback and write its own
# /run + /var/lib/net-safe state — subtle breakage there is worse than the gain). What we DO take
# is the cheap, can't-break-it surface: cap the ceiling to exactly what it uses (nft/ip/conntrack
# = CAP_NET_ADMIN, ping = CAP_NET_RAW; everything else dropped) and the non-filesystem protections.
TimeoutStopSec=15
NoNewPrivileges=true
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
ProtectHome=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectClock=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
