# bastion/templates/systemd/edge-ai-collect.service
[Unit]
Description=edge-ai signal collector (root, READ-ONLY, sanitizes -> external-safe signals.json)
Documentation=file:/usr/local/sbin/edge-ai-collect
After=crowdsec.service
Wants=crowdsec.service

[Service]
Type=oneshot
# Supplies NFT_TABLE/MODE so the collector reads the right managed table (inet edge vs
# inet bastion). `-` => optional: a not-yet-generated env file is not a startup failure.
EnvironmentFile=-/etc/bastion/machine.env
ExecStart=/usr/local/sbin/edge-ai-collect
# Privileged read-only collector (reads cscli + journal). Holds ZERO firewall-write power.
ProtectSystem=strict
ReadWritePaths=/var/lib/edge-ai
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
