# bastion/templates/systemd/edge-feed.timer
[Unit]
Description=daily edge threat-feed refresh

[Timer]
OnBootSec=5min
OnCalendar=daily
RandomizedDelaySec=30min
Persistent=true

[Install]
WantedBy=timers.target
